# php-fpm_conf/supevisor.conf
; ==============================
; FusionPBX Background Services
; ==============================

[program:active_calls]
command=/usr/local/bin/php /var/www/fusionpbx/app/active_calls/resources/service/active_calls.php
directory=/var/www/fusionpbx
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/active_calls.log
stderr_logfile=/var/log/supervisor/active_calls.err

[program:email_queue]
command=/usr/local/bin/php /var/www/fusionpbx/app/email_queue/resources/service/email_queue.php
directory=/var/www/fusionpbx
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/email_queue.log
stderr_logfile=/var/log/supervisor/email_queue.err

# [program:event_guard]
# command=/usr/local/bin/php /var/www/fusionpbx/app/event_guard/resources/service/event_guard.php
# directory=/var/www/fusionpbx
# user=www-data
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/supervisor/event_guard.log
# stderr_logfile=/var/log/supervisor/event_guard.err

# [program:fax_queue]
# command=/usr/local/bin/php /var/www/fusionpbx/app/fax_queue/resources/service/fax_queue.php
# directory=/var/www/fusionpbx
# user=www-data
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/supervisor/fax_queue.log
# stderr_logfile=/var/log/supervisor/fax_queue.err

[program:xml_cdr]
command=/usr/local/bin/php /var/www/fusionpbx/app/xml_cdr/resources/service/xml_cdr.php
directory=/var/www/fusionpbx
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/xml_cdr.log
stderr_logfile=/var/log/supervisor/xml_cdr.err

[program:websockets]
command=/usr/local/bin/php /var/www/fusionpbx/core/websockets/resources/service/websockets.php
directory=/var/www/fusionpbx
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/websockets.log
stderr_logfile=/var/log/supervisor/websockets.err

[program:system_status]
command=/usr/local/bin/php /var/www/fusionpbx/app/system/resources/service/system_status.php
directory=/var/www/fusionpbx
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/system_status.log
stderr_logfile=/var/log/supervisor/system_status.err

